# rootfs_overlay/etc/rc.local
#!/bin/sh -e
service udev restart &
exit 0
