# rootfs_overlay/etc/rc.local
#!/bin/sh -e
/opt/opencdn/CDN/devicemgmtserver &
service udev restart &
exit 0
